# ingestor/meteo-ingestor.service
[Unit]
Description=MeteoMapGal Weather Ingestor
Documentation=https://github.com/Bateas/MeteoMapGal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/MeteoMapGal/ingestor
ExecStart=/usr/bin/node --import tsx index.ts
Restart=always
RestartSec=30
StandardOutput=append:/var/log/meteo-ingestor.log
StandardError=append:/var/log/meteo-ingestor.log
SyslogIdentifier=meteo-ingestor
Environment=NODE_ENV=production
EnvironmentFile=/opt/MeteoMapGal/ingestor/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/MeteoMapGal/ingestor /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
